# ops/nag-deploy.service
[Unit]
Description=Pull latest nag from GitHub and restart the Expo dev server

[Service]
Type=oneshot
WorkingDirectory=%h/nag
EnvironmentFile=-%h/.config/nag/deploy.env
ExecStart=%h/nag/ops/deploy.sh
