# ops/nag-expo.service
[Unit]
Description=nag Expo dev server (Metro bundler)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/nag
EnvironmentFile=-%h/.config/nag/expo.env
# pnpm is expected to be on $PATH (install via corepack, mise, nvm, or volta).
# --dev-client so the custom dev-client build connects (not Expo Go).
ExecStart=/usr/bin/env pnpm --filter @nag/app exec expo start --dev-client --port 8081
Restart=on-failure
RestartSec=5
# Metro uses a lot of watchers; bump the per-user limit via sysctl separately
# (see docs/DevServerHosting.md).

[Install]
WantedBy=default.target
